# pkg/regin_0.1.0_amd64/DEBIAN/postinst
#!/bin/sh
set -e

# Create regin system user if it doesn't exist
if ! getent passwd regin >/dev/null 2>&1; then
    adduser --system --group --home /var/lib/regin --no-create-home \
        --gecos "Regin monitoring daemon" regin
fi

# Create state directories
mkdir -p /var/lib/regin/skills
chown -R regin:regin /var/lib/regin

# Protect config (contains API key) — readable by regin group only
chown root:regin /etc/regin/config.toml
chmod 640 /etc/regin/config.toml

# If user hasn't customised config yet, seed the skills directory
# from shipped examples
if [ ! -d /var/lib/regin/skills/disk-usage ]; then
    cp -r /usr/share/regin/skills/* /var/lib/regin/skills/
    chown -R regin:regin /var/lib/regin/skills
fi

# Enable and start the service
if [ -d /run/systemd/system ]; then
    systemctl daemon-reload
    systemctl enable regind.service
    # Don't start automatically — user needs to set API key first
    echo "  Regin installed. Set your NanoGPT API key in /etc/regin/config.toml"
    echo "  then start with: sudo systemctl start regind"
fi

#DEBHELPER#
exit 0
