# pkg/regin_0.1.0_amd64/DEBIAN/prerm
#!/bin/sh
set -e

if [ -d /run/systemd/system ]; then
    systemctl stop regind.service 2>/dev/null || true
    systemctl disable regind.service 2>/dev/null || true
fi

#DEBHELPER#
exit 0
