# pkg/regin_0.1.0_amd64/usr/lib/systemd/system/regind.service
[Unit]
Description=Regin Monitoring Daemon
Documentation=https://github.com/exedev/regin
After=network.target

[Service]
Type=exec
User=regin
Group=regin
WorkingDirectory=/var/lib/regin
ExecStart=/usr/bin/regind --config /etc/regin/config.toml
Restart=on-failure
RestartSec=10

# Hardening
ProtectSystem=strict
ReadWritePaths=/var/lib/regin
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

Environment=RUST_LOG=info
StandardOutput=journal
StandardError=journal
SyslogIdentifier=regind

[Install]
WantedBy=multi-user.target
